# rootfs/etc/services.d/ipsec/run
#!/usr/bin/execlineb -P

ipsec start --nofork
